# matcha_config.conf
description "uWSGI server incstance cocnfigured to serve matcha"

start on runlevel [2345]
stop on runlevel [!2345]

setuid vagrant
setgid www-data

env PATH=/vagrant/app_venv/bin
chdir /vagrant
exec uwsgi --ini uwsgi_config.ini
